# roles/backend/templates/systemd.service
[Unit]
Description = {{COMPONENT}} Service

[Service]
User={{USER}}
ExecStart=/usr/bin/node index.js
SyslogIdentifier={{COMPONENT}}
WorkingDirectory=/home/{{USER}}/{{COMPONENT}}/
Environment="MONGO_ENDPOINT=mongodb+srv://{{DB_USER}}:{{DB_PASS}}@{{DB_URL}}/login-app-db?retryWrites=true&w=majority"
Environment="APP_VERSION={{APP_VERSION}}"
Environment="ENV={{ENV}}"
  # StandardOutput=syslog
  # StandardError=syslog

[Install]
WantedBy=multi-user.target